# service/wallpaper.service
[Unit]
Description=Change Wallpaper

[Service]
ExecStart=/usr/local/bin/swww_img_changer
